# launchd/com.rust-sci-hub-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier - must be unique on the system -->
    <key>Label</key>
    <string>com.rust-sci-hub-mcp</string>
    
    <!-- Program to execute -->
    <!-- Note: Paths with ~ will be replaced during installation -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/rust-sci-hub-mcp</string>
        <string>--daemon</string>
        <string>--config</string>
        <string>HOME_DIR/.config/rust-sci-hub-mcp/config.toml</string>
        <string>--pid-file</string>
        <string>HOME_DIR/Library/Application Support/rust-sci-hub-mcp/rust-sci-hub-mcp.pid</string>
    </array>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>HOME_DIR/Library/Application Support/rust-sci-hub-mcp</string>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>HOME_DIR</string>
    </dict>
    
    <!-- Run at login -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Keep the service alive -->
    <key>KeepAlive</key>
    <dict>
        <!-- Restart if it crashes -->
        <key>Crashed</key>
        <true/>
        
        <!-- Don't restart if it exits successfully -->
        <key>SuccessfulExit</key>
        <false/>
        
        <!-- Network state requirements -->
        <key>NetworkState</key>
        <true/>
    </dict>
    
    <!-- Restart throttling -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <!-- Process type -->
    <key>ProcessType</key>
    <string>Background</string>
    
    <!-- Nice value (priority) -->
    <key>Nice</key>
    <integer>5</integer>
    
    <!-- Resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <!-- Max number of open files -->
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
    
    <key>HardResourceLimits</key>
    <dict>
        <!-- Hard limit for open files -->
        <key>NumberOfFiles</key>
        <integer>8192</integer>
    </dict>
    
    <!-- Standard output/error logging -->
    <key>StandardOutPath</key>
    <string>HOME_DIR/Library/Logs/rust-sci-hub-mcp/stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>HOME_DIR/Library/Logs/rust-sci-hub-mcp/stderr.log</string>
    
    <!-- Exit timeout (seconds) - time to wait for graceful shutdown -->
    <key>ExitTimeOut</key>
    <integer>30</integer>
    
    <!-- Enable transactions (for macOS 10.14+) -->
    <key>EnableTransactions</key>
    <true/>
    
    <!-- Enable pressured exit (for memory pressure events) -->
    <key>EnablePressuredExit</key>
    <true/>
</dict>
</plist>